# templates/agent/guestkit-agent.service
[Unit]
Description=GuestKit Agent (Zyvor)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/bin/sh -c 'until [ -e /dev/virtio-ports/org.qemu.guest_agent.0 ]; do sleep 2; done; exec /usr/local/bin/guestkit agent --channel virtio'
Restart=always
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
