# templates/agent/zyvor-guest-agent.service
[Unit]
Description=Zyvor VM Tools Guest Agent
Documentation=https://zyvor.dev/guestkit
After=network-online.target
Wants=network-online.target
ConditionPathExists=/dev/virtio-ports/org.qemu.guest_agent.0

[Service]
Type=simple
ExecStart=/usr/bin/zyvor-guest-agent
Restart=always
RestartSec=5
User=root
StandardOutput=journal
StandardError=journal
Environment=RUST_LOG=info
RuntimeDirectory=zyvor
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
