# templates/agent/zyvor-guest-updater.service
[Unit]
Description=Zyvor GuestAgent signed update check
Documentation=https://zyvor.dev/guestkit

[Service]
Type=oneshot
ExecStart=/usr/bin/zyvor-guest-agent --scheduled-update
User=zyvor-agent
Group=zyvor-agent

[Install]
WantedBy=multi-user.target
